# deploy/polymarket-dashboard.service
[Unit]
Description=Polymarket BTC Bot Dashboard (Flask)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=alex
WorkingDirectory=/home/alex/polymarket-btc-up-down
ExecStart=/home/alex/polymarket-btc-up-down/venv/bin/python dashboard.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONIOENCODING=utf-8

[Install]
WantedBy=multi-user.target
